# systemd/home-manager-switch.service
[Unit]
Description=Apply home-manager configuration on startup
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
ExecStart=nix run nixpkgs#home-manager -- switch --flake /home/coder/dotfiles#coder
RemainAfterExit=yes
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=default.target
